# snowfall.openrc
#!/sbin/openrc-run
# snowfall — YetiOS graphical login manager

name="snowfall"
description="Snowfall login manager"
command="/usr/local/sbin/snowfall"
command_background=true
pidfile="/run/snowfall.pid"
output_log="/var/log/snowfall.log"
error_log="/var/log/snowfall.log"

depend() {
    need udev
    after snowcone
    before local
}

start_pre() {
    checkpath -d /usr/share/wayland-sessions
    checkpath -f -m 0644 /var/log/snowfall.log

    # Stop the boot splash so it releases DRM master.
    if rc-service -e snowcone 2>/dev/null; then
        if rc-service snowcone status 2>/dev/null | grep -q started; then
            ebegin "Stopping snowcone splash"
            rc-service snowcone stop 2>/dev/null
            sleep 1
            eend 0
        fi
    fi
}
